# dockerfile/mysql/config_file/mysql.ini
[program:mysqld]
command=/etc/init.d/mysql.server start
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/tmp/supervisor_mysql.log
