# dockerfile/cachecloud/cachec/cachecloud.ini
[program:cachecloud]
command=/bin/sh /opt/cachecloud-web/start.sh
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/tmp/supervisor_mysql.log
