# conf/monitor_prometheus/consul/systemd-service/iprinit.service
[Unit]
Description=IBM Power Raid init daemon

[Service]
Type=forking
ExecStart=/usr/sbin/iprinit --daemon

[Install]
WantedBy=iprutils.target
