# conf/monitor_prometheus/consul/systemd-service/iprupdate.service
[Unit]
Description=IBM Power Raid update daemon

[Service]
Type=forking
ExecStart=/usr/sbin/iprupdate --daemon

[Install]
WantedBy=iprutils.target
